# app.service
[Unit]
Description=Affiche Hello KONDI CODE
After=network.target

[Service]
WorkingDirectory=/var/www/ua-app
ExecStart=/var/www/ua-app/venv/bin/gunicorn ua/wsgi:application

Type=simple
User=root
Group=root

[Install]
WantedBy=default.target